# src/OwlClient/owl.service
[Unit]
After=network.target

[Service]
ExecStart=/home/nick/HomeMonitor/src/OwlClient/startup.sh
Restart=always

[Install]
WantedBy=default.target
