# src/RakoClient/rako.service
[Unit]
After=network.target

[Service]
ExecStart=/home/nick/HomeMonitor/src/RakoClient/startup.sh

[Install]
WantedBy=default.target
